# scripts/clarateach-agent.service
[Unit]
Description=ClaraTeach Worker Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=IMAGES_DIR=/var/lib/clarateach/images
Environment=SOCKET_DIR=/tmp/clarateach
ExecStartPre=/bin/bash -c '\
  METADATA_URL="http://metadata.google.internal/computeMetadata/v1/instance/attributes"; \
  METADATA_HEADER="Metadata-Flavor: Google"; \
  echo "" > /run/clarateach-agent.env; \
  RESPONSE=$(curl -s -w "%%{http_code}" -o /tmp/meta-val.txt -H "$METADATA_HEADER" "$METADATA_URL/agent-token" 2>/dev/null); \
  if [ "$RESPONSE" = "200" ]; then echo "AGENT_TOKEN=$(cat /tmp/meta-val.txt)" >> /run/clarateach-agent.env; fi; \
  RESPONSE=$(curl -s -w "%%{http_code}" -o /tmp/meta-val.txt -H "$METADATA_HEADER" "$METADATA_URL/workspace-token-secret" 2>/dev/null); \
  if [ "$RESPONSE" = "200" ]; then echo "WORKSPACE_TOKEN_SECRET=$(cat /tmp/meta-val.txt)" >> /run/clarateach-agent.env; fi; \
  rm -f /tmp/meta-val.txt'
EnvironmentFile=-/run/clarateach-agent.env
ExecStart=/usr/local/bin/agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
